# roles/loki/templates/loki.service.j2
# Managed by Ansible — do not edit by hand
[Unit]
Description=Loki log aggregation system
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ loki_user }}
Group={{ loki_user }}
ExecStart={{ loki_install_dir }}/loki \
  -config.file={{ loki_config_path }}/loki-config.yml
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ loki_storage_path }}

[Install]
WantedBy=multi-user.target
